# com.cpufreestyle.middle-button-typer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cpufreestyle.middle-button-typer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/CHANGE_ME/Applications/MiddleButtonTyper.app/Contents/MacOS/MiddleButtonTyper</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/CHANGE_ME</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/CHANGE_ME/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/middle_button_typer.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/middle_button_typer.log</string>
</dict>
</plist>
