# com.a1-6.middle-button-typer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.a1-6.middle-button-typer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/a1-6/.qclaw/workspace/middle-button-typer/MiddleButtonTyper</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/middle_button_typer.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/middle_button_typer.log</string>
</dict>
</plist>
